# scripts/zoplog-dns-dispatcher.service
# ZopLog DNS Dispatcher Service
# This service ensures that /etc/resolv.conf is automatically updated with the correct
# DNS server (the current default gateway IP) on system boot and network changes.
# This is crucial for DNS resolution to work properly on the ZopLog system itself,
# especially when the internet gateway IP changes.

[Unit]
Description=ZopLog DNS Dispatcher
After=network.target
Wants=network.target

[Service]
Type=oneshot
User=zoplog
Group=zoplog
ExecStart=/etc/NetworkManager/dispatcher.d/zoplog-dns-dispatcher
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
